# scripts/save/choggiung.wifi.conf
start on (filesystem
          and static-network-up
          and runlevel [2345])
stop on runlevel [016]

exec /bin/bash /home/administrator/scripts/wifi.sh
respawn
