# install/etc/init/innproxy.usagewatch.conf
start on (filesystem
          and static-network-up
          and runlevel [2345])
stop on runlevel [016]

exec /bin/bash /usr/local/bin/usagewatch
respawn
